# etc/supervisor.d/echo.conf
[program:echo]
process_name=%(program_name)s_%(process_num)02d
command=/opt/homebrew/bin/node /opt/homebrew/bin/npx laravel-echo-server start
directory=/Users/alex/code/processmaker/
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=alex
numprocs=1
redirect_stderr=true
stdout_logfile=/Users/alex/logs/%(program_name)s_%(process_num)02d.log
stopwaitsecs=3600
